# openresty/nginx/nginx.cron
0 * *  *  *  root /bin/sh /data/scripts/nginx_log.sh > /dev/null 2>&1 &
