# nginx/nginx.cron
0 * *  *  *  root /bin/sh /data/nginx/nginx_log.sh > /dev/null 2>&1 &
